# cookbooks/dns/files/duck-dns.service
[Unit]
Description=Duck DNS update

[Service]
Type=oneshot

ExecStart=/var/duck-dns/duck.sh
